# software/main/pcf.service
[Unit]
Description=pcf.service Service

After=network.target

[Service]
ExecStart=/usr/bin/python /home/pi/main_pcf_ver4.py
WorkingDirectory=/home/pi/S
StandardOutput=inherit
StandardError=inherit
Restart=always
User=root

[Install]
WantedBy=multi-user.target